# lldpad-post
#!/sbin/openrc-run
LTOOL='/usr/sbin/lldptool'
depend() {
    need lldpad
    after lldpad
}

start() {
    for net in `ls /sys/class/net/ | grep 'et\|en'`; do
        echo "Enabling LLDP on $net"
        $LTOOL set-lldp -i $net adminStatus=rxtx >/dev/null
        for tlv in `lldptool --h |grep -A100 "TLV identifiers:" | tail +2 | awk -F':' '{print $1}' | xargs`; do
            $LTOOL -T -i $net -V $tlv enableTx=yes >/dev/null
        done
    done
    true
}

stop() {
    for net in `ls /sys/class/net/ | grep 'et\|en'`; do
        echo "Disabling LLDP on $net"
        $LTOOL set-lldp -i $net adminStatus=disabled >/dev/null
    done
    true
}
